# systemd/pape.service
[Unit]
Description=random pape
[Service]
Type=oneshot
ExecStart=%h/.local/bin/pape ?
