# systemd/pape.timer
[Unit]
Description=set pape
[Timer]
OnCalendar=*:0/15
Persistent=true
[Install]
WantedBy=default.target
